# deploy/rune.service
[Unit]
Description=RUNE Web REPL
After=local-fs.target

[Service]
Type=simple
User=@@SERVICE_USER@@
Group=@@PROXY_GROUP@@
WorkingDirectory=@@APP_DIR@@
RuntimeDirectory=@@SERVICE_NAME@@
RuntimeDirectoryMode=0750

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# The socket directory is accessible only to the runtime user and proxy group,
# so every possible peer is trusted to provide Caddy's forwarding headers.
ExecStart=@@APP_DIR@@/.venv/bin/python -m uvicorn app:app --app-dir @@APP_DIR@@/web --uds /run/@@SERVICE_NAME@@/rune.sock --workers 1 --proxy-headers --forwarded-allow-ips=* --no-server-header --limit-concurrency 32 --backlog 64 --timeout-keep-alive 5

Restart=on-failure
RestartSec=3
TimeoutStopSec=10
KillMode=control-group

NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
PrivateNetwork=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true
UMask=0007

# Give each service start a fresh, size-bounded, non-executable scratch area.
TemporaryFileSystem=/tmp:rw,nosuid,nodev,noexec,size=32M,mode=1777
TemporaryFileSystem=/var/tmp:rw,nosuid,nodev,noexec,size=8M,mode=1777

# A future interpreter regression must not be able to fill the host journal.
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

# Tune these for the VPS if necessary. The cgroup limits include disposable
# evaluator subprocesses, not just the Uvicorn parent.
MemoryHigh=448M
MemoryMax=512M
MemorySwapMax=0
TasksMax=64
CPUQuota=100%
LimitNOFILE=256

[Install]
WantedBy=multi-user.target
